# scripts/watchdog.cron
*/5 * * * * root cd /var/hivebot; scripts/watchdog.rb > /var/log/somelogfile.log 2>&1;
